# scripts/backup.cron
# Run backup.sh at 2 AM every day
0 2 * * * /corpus/backup.sh >> /var/log/cron.log 2>&1